# mac/essentials/Cellframe-Dashboard-pkginstall/postinstall
#!/bin/bash

APP_NAME=Cellframe-Dashboard
DAP_PREFIX=/Users/$USER/Applications/Cellframe.app/Contents/Resources
set -x
sudo xattr -rd com.apple.quarantine /Applications/"$APP_NAME".app
sudo killall -HUP mDNSResponder

EXEC_PREFIX=/Applications/$APP_NAME.app/Contents/MacOS
SERVICE=com.demlabs."$APP_NAME"Service
NODE=com.demlabs.cellframe-node

mkdir -p $DAP_PREFIX
sudo mv -f /Applications/$APP_NAME.app/Contents/Resources/etc $DAP_PREFIX
sudo mv -f /Applications/$APP_NAME.app/Contents/Resources/share $DAP_PREFIX

for filename in $(find $DAP_PREFIX | grep -v bugreport); do
    if [ -d $filename ]; then
        chmod 0775 $filename || true
    else
        chmod 0664 $filename || true
    fi
done

wd=$(pwd)
cd $EXEC_PREFIX

sudo chmod +x $EXEC_PREFIX/create_configs.sh
sudo chmod +x $EXEC_PREFIX/create_configs_from_tpl.sh

$EXEC_PREFIX/create_configs.sh 2>> /tmp/cellframe-dashboard_error.txt || { error=$?; echo "not create configs $error" >> /tmp/cellframe-dashboard_error.txt; exit $error; }

rm $EXEC_PREFIX/*.sh

for filename in $(ls . | grep -v '.'); do
	if [ -L $filename ] && [ ! -e $filename ]; then
		strip -u -r $filename
	fi
done

cd $wd

echo $SERVICE
sudo launchctl stop $SERVICE
sudo launchctl unload -w ~/Library/LaunchAgents/$SERVICE.plist

echo $NODE || exit $?
sudo launchctl stop $NODE
sudo launchctl unload -w ~/Library/LaunchAgents/$NODE

ln -sf /Applications/"$APP_NAME".app/Contents/Resources/$SERVICE.plist ~/Library/LaunchAgents/$SERVICE.plist
ln -sf /Applications/"$APP_NAME".app/Contents/Resources/$NODE.plist ~/Library/LaunchAgents/$NODE.plist

launchctl bootstrap gui/$UID ~/Library/LaunchAgents/$NODE.plist
echo $UID > /tmp/my_uid
launchctl load -w ~/Library/LaunchAgents/$SERVICE.plist


launchctl start $NODE
launchctl start $SERVICE

